# kubernetes/xrootd/build/cert-renew/files/etc/supervisord.d/10-certrenew.conf
[program:certrenew]
command=bash -c "cd / && /renewal.sh"
autorestart=true
