# kubernetes/xrootd/build/xrootd-mon/config/default/etc/supervisord.d/10-xrootd.conf
[program:xrootd-mon]
command=python3 /opt/xroot-mon.py
autorestart=true

[program:gunicorn]
command=bash -c "cd /opt; gunicorn --workers=2  http-server:app -b 0.0.0.0:8000"
autorestart=true
